# server/sample.service
[Unit]
Description=<Name Here>

[Service]
ExecStart=<path to yarn> start
Restart=always
RestartSec=60
WorkingDirectory=<directory of the folder>
